# rpi-status-display.service
[Unit]
Description=rpi-status-display

[Service]
ExecStart=/usr/bin/rpi-status-display

[Install]
WantedBy=multi-user.target
